# container/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:phpd]
command=php-fpm
